# fr24_monitor.cron
# FR24 Monitor Crontab
# Runs the FR24 monitoring script every 10 minutes
# Format: minute hour day month weekday command

# Run FR24 monitor every 10 minutes
*/10 * * * * /home/john/fr24-mon/fr24_monitor.sh --log-file /home/john/fr24-mon/fr24_monitor.log

# Alternative with more verbose logging (uncomment if needed)
# */10 * * * * /home/john/fr24-mon/fr24_monitor.sh --verbose --log-file /home/john/fr24-mon/fr24_monitor.log

# Alternative with custom settings (examples - uncomment and modify as needed)
# */10 * * * * /home/john/fr24-mon/fr24_monitor.sh --threshold 5 --min-uptime 1 --log-file /home/john/fr24-mon/fr24_monitor.log
# */10 * * * * /home/john/fr24-mon/fr24_monitor.sh --service-name piaware --log-file /home/john/fr24-mon/fr24_monitor.log

# Run a more frequent check during peak hours (every 5 minutes between 6 AM and 11 PM)
# */5 6-23 * * * /home/john/fr24-mon/fr24_monitor.sh --log-file /home/john/fr24-mon/fr24_monitor.log

# Run less frequent check during night hours (every 30 minutes between midnight and 6 AM)
# */30 0-5 * * * /home/john/fr24-mon/fr24_monitor.sh --log-file /home/john/fr24-mon/fr24_monitor.log
